# HelloString/systemctl/helbus.service
[Unit]
Description=hel dbus...

[Service]
ExecStart=/usr/bin/dbus-daemon --session --nofork --address="unix:path=/var/run/dbus/hel_bus"
